# autostart/cfw.desktop
[Desktop Entry]
Type=Application
NoDisplay=true
Terminal=false
Name=Clash
Exec=$HOME/Apps/clash/cfw
