# autostart/goldendict.desktop
[Desktop Entry]
Type=Application
NoDisplay=true
Terminal=false
Name=GoldenDict
Exec=/usr/bin/goldendict

